# etc/supervisord.d/slurmd.ini
[program:slurmd]
command=/usr/local/sbin/slurmd -D
autostart=true
redirect_stderr=true
#stdout_syslog=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
